# services/hydra/run
#!/bin/bash
cd /app/hydra
/sbin/setuser abc python nzbhydra.py --nobrowser --config /config
