# services/hydra/run
#!/bin/bash
/sbin/setuser abc python /config/hydra/nzbhydra.py --nobrowser
