# services/hydra/run
#!/bin/bash
/sbin/setuser abc python nzbhydra.py --nobrowser
