# services/hydra/run
#!/bin/bash
cd /config/hydra
/sbin/setuser abc python nzbhydra.py --nobrowser
