# root/etc/services.d/hydra/run
#!/usr/bin/with-contenv bash

cd /config/hydra || exit

exec \
	s6-setuidgid abc python /config/hydra/nzbhydra.py \
	--nobrowser
