# root/etc/services.d/hydra/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid abc python /config/hydra/nzbhydra.py --nobrowser
